# aws/ami/stig/roles/canonical-ubuntu-18.04-lts-stig-hardening/files/klog-log-permissions.path
[Unit]
Description=Monitor klog paths
[Path]
PathChanged=/var/log/kubernetes/
PathChanged=/var/log/antrea/
Unit=klog-fix.service
[Install]
WantedBy=multi-user.target